# units/user-configdrive.path
[Unit]
Description=Watch for a cloud-config at /media/configdrive

[Path]
DirectoryNotEmpty=/media/configdrive
