# units/user-cloudinit@.path
[Unit]
Description=Watch for a cloud-config at %f

# Inhibit all cloudinit-related activities if ignition disabled us.
ConditionPathExists=!/etc/.disable-cloudinit

[Path]
PathExists=%f
